# dashboard.service
[Unit]
Description=Grafana-like Dashboard
Documentation=https://github.com/yourusername/your-repo
After=network.target

[Service]
Type=simple
User=dashboard
Environment=NODE_ENV=production
WorkingDirectory=/opt/dashboard
ExecStart=/usr/bin/node src/server.js --config /etc/dashboard/config.ini
Restart=always
RestartSec=10s
TimeoutStopSec=30s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target 